# systemd/.config/systemd/user/battery-notifier.service
[Unit]
Description=Battery notifier script

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-notifier.zsh

